# supervisor.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:api]
command=/app/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/api.log
stderr_logfile=/app/logs/api_error.log
environment=PYTHONPATH="/app"

[program:celery-worker]
command=/app/.venv/bin/celery -A scheduler.celery_app worker --loglevel=info --concurrency=4
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/celery_worker.log
stderr_logfile=/app/logs/celery_worker_error.log
environment=PYTHONPATH="/app"

[program:celery-beat]
command=/app/.venv/bin/celery -A scheduler.celery_app beat --loglevel=info
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/celery_beat.log
stderr_logfile=/app/logs/celery_beat_error.log
environment=PYTHONPATH="/app"

[program:telegram-bot]
command=/app/.venv/bin/python -c "from integrations.telegram_bot import run_bot; run_bot()"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/telegram_bot.log
stderr_logfile=/app/logs/telegram_bot_error.log
environment=PYTHONPATH="/app"

[program:dashboard]
command=/app/.venv/bin/streamlit run monitoring/dashboard.py --server.port 8501 --server.address 0.0.0.0
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/dashboard.log
stderr_logfile=/app/logs/dashboard_error.log
environment=PYTHONPATH="/app"
